# evidence-s3-sync.service
[Unit]
Description=Evidence Files S3 Sync Service
After=network.target

[Service]
Type=oneshot
User=ec2-user
WorkingDirectory=/opt/msp-checklist-system
ExecStart=/opt/msp-checklist-system/evidence-s3-sync.sh
Environment=EVIDENCE_STORAGE_PATH=/opt/msp-checklist-system/evidence-files
Environment=EVIDENCE_S3_BUCKET=your-bucket-name
Environment=EVIDENCE_S3_PREFIX=evidence-files

[Install]
WantedBy=multi-user.target
